# deploy/systemd/openjarvis-healthcheck.service
[Unit]
Description=OpenJarvis health probe
After=network.target

[Service]
Type=oneshot
User=openjarvis
WorkingDirectory=/opt/openjarvis
ExecStart=/opt/openjarvis/scripts/check_openjarvis_stack.sh
Environment=OPENJARVIS_BASE_URL=https://jarvis.example.com
